# systemd/sing-box-tproxy.service
[Unit]
Description=sing-box TProxy policy routing and nftables rules
Documentation=https://sing-box.sagernet.org
After=network-online.target
Wants=network-online.target
Before=sing-box.service

[Service]
Type=oneshot
# IPv6 PD/RA 前缀可能随运营商变化；每次应用 nft 前都按当前网卡状态重生成脚本，避免使用过期前缀。
ExecStartPre=/usr/local/sbin/refresh-sing-box-tproxy-setup
ExecStart=/usr/local/sbin/sing-box-tproxy-setup
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
